# services/axon-sandbox/axon-sandbox.service
[Unit]
Description=Axon OS Sandbox Manager Service
After=dbus.socket
Requires=dbus.socket

[Service]
Type=dbus
BusName=org.axonos.Sandbox
ExecStart=/usr/bin/python3 AXON_SERVICES_DIR/axon-sandbox/sandbox_manager.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target
